# systemd/myos-reconcile.service
[Unit]
Description=Reconcile myOS operator accounts

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/reconcile-operators.py
